# docker/root/etc/services.d/php-fpm/run
#!/command/execlineb -P
s6-setuidgid abc

php-fpm8 -F
