# docker/root/etc/services.d/laravel-worker/run
#!/command/execlineb -P
s6-setuidgid abc

php /app/artisan queue:work
