# docker/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
with-contenv
s6-setuidgid abc

nginx -g "pid /tmp/nginx.pid; daemon off;"
