# docker/root/etc/services.d/nginx/run
#!/usr/bin/execlineb -P
s6-setuidgid abc

nginx -g "pid /tmp/nginx.pid; daemon off;"
