# docker/root/etc/services.d/laravel-websockets/run
#!/usr/bin/execlineb -P
s6-setuidgid abc

php artisan websockets:serve
